# avvedrp/trunk/StudentActivityVisualization/cronjob/crontab
* * * * * root /root/Kettle-3.0.2/mdl >>/root/dive/txform4.log 2>&1
